# 70-g610.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c333", ACTION=="add", RUN+="/bin/sh -c 'DISPLAY=:0 /usr/local/bin/g610-tool'"
